# docs/systemd/estimator-cleanup.service
[Unit]
Description=AI Estimator System - Auto Cleanup Batch (GDPR Compliance)
After=network.target

[Service]
Type=oneshot
User=your-username
WorkingDirectory=/path/to/ai-estimator2/backend
Environment="PATH=/path/to/python/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/path/to/python/bin/python -m app.tasks.cleanup
StandardOutput=journal
StandardError=journal
SyslogIdentifier=estimator-cleanup

[Install]
WantedBy=multi-user.target
